# oled-i2c-screen-adapter/rootfs/etc/services.d/display_controller/run
#!/bin/sh
CONFIG_FILE="/data/options.json"

# MQTT configuration
MQTT_BROKER=$(jq -r '.mqtt_broker' "$CONFIG_FILE")
MQTT_PORT=$(jq -r '.mqtt_port' "$CONFIG_FILE")
MQTT_USER=$(jq -r '.mqtt_user' "$CONFIG_FILE")
MQTT_PASSWORD=$(jq -r '.mqtt_password' "$CONFIG_FILE")

# I2C and display configuration
I2C_ADDRESS=$(jq -r '.i2c_address' "$CONFIG_FILE")
I2C_PORT=$(jq -r '.i2c_port' "$CONFIG_FILE")
DISPLAY_WIDTH=$(jq -r '.display_width' "$CONFIG_FILE")
DISPLAY_HEIGHT=$(jq -r '.display_height' "$CONFIG_FILE")
DISPLAY_TYPE=$(jq -r '.display_type' "$CONFIG_FILE")

# Display parameters configuration
REFRESH_INTERVAL=$(jq -r '.refresh_interval' "$CONFIG_FILE")
DEFAULT_BRIGHTNESS=$(jq -r '.default_brightness' "$CONFIG_FILE")
FONT_SIZE=$(jq -r '.font_size' "$CONFIG_FILE")
SHOW_HEADER=$(jq -r '.show_header' "$CONFIG_FILE")
CUSTOM_LINES=$(jq -c '.custom_lines' "$CONFIG_FILE")

# MQTT topics configuration
MQTT_TOPIC_TEXT=$(jq -r '.topics.text' "$CONFIG_FILE")
MQTT_TOPIC_MODE=$(jq -r '.topics.mode' "$CONFIG_FILE")
MQTT_TOPIC_BRIGHTNESS=$(jq -r '.topics.brightness' "$CONFIG_FILE")
MQTT_TOPIC_REFRESH=$(jq -r '.topics.refresh' "$CONFIG_FILE")
MQTT_TOPIC_STATUS=$(jq -r '.topics.status' "$CONFIG_FILE")

# Misc configuration
QR_LINK=$(jq -r '.qr_link' "$CONFIG_FILE")
LOG_LEVEL=$(jq -r '.log_level' "$CONFIG_FILE")

# Export environment variables
export MQTT_BROKER
export MQTT_PORT
export MQTT_USER
export MQTT_PASSWORD
export I2C_ADDRESS
export I2C_PORT
export DISPLAY_WIDTH
export DISPLAY_HEIGHT
export DISPLAY_TYPE
export REFRESH_INTERVAL
export DEFAULT_BRIGHTNESS
export MQTT_TOPIC_TEXT
export MQTT_TOPIC_COMMAND
export MQTT_TOPIC_MODE
export MQTT_TOPIC_BRIGHTNESS
export MQTT_TOPIC_REFRESH
export QR_LINK
export LOG_LEVEL
export FONT_SIZE
export SHOW_HEADER
export CUSTOM_LINES
export MQTT_TOPIC_STATUS



# Logging helpers (respect LOG_LEVEL)
LOG_LEVEL=${LOG_LEVEL:-INFO}
level_to_num() {
    case "$1" in
        DEBUG) echo 10 ;;
        INFO) echo 20 ;;
        WARNING) echo 30 ;;
        ERROR) echo 40 ;;
        CRITICAL) echo 50 ;;
        *) echo 20 ;;
    esac
}
log() {
    lvl="$1"; shift
    msg="$*"
    [ "$(level_to_num "$lvl")" -ge "$(level_to_num "$LOG_LEVEL")" ] || return 0
    printf '[%s] %s: %s\n' "$(date '+%Y-%m-%d %H:%M:%S')" "$lvl" "$msg"
}
log_debug() { log DEBUG "$@"; }
log_info() { log INFO "$@"; }
log_warning() { log WARNING "$@"; }
log_error() { log ERROR "$@"; }
log_critical() { log CRITICAL "$@"; }

log_info "Configuration loaded"
log_info "MQTT broker: $MQTT_BROKER:$MQTT_PORT"
log_info "MQTT user: ${MQTT_USER:-(none)}"
log_info "I2C config: port $I2C_PORT, address $I2C_ADDRESS"
log_info "Display config: ${DISPLAY_WIDTH}x${DISPLAY_HEIGHT} ($DISPLAY_TYPE)"
log_info "Refresh: ${REFRESH_INTERVAL}s, brightness: $DEFAULT_BRIGHTNESS"
log_info "MQTT topics: text=$MQTT_TOPIC_TEXT, command=$MQTT_TOPIC_COMMAND, mode=$MQTT_TOPIC_MODE, brightness=$MQTT_TOPIC_BRIGHTNESS, refresh=$MQTT_TOPIC_REFRESH"

# Check Python dependencies
log_info "Checking Python dependencies..."
python3 -c "import os, sys; print(f'Python version: {sys.version}')"
python3 -c "import json; print('JSON: OK')" 2>/dev/null || log_error "Missing module: json"
python3 -c "import paho.mqtt.client as mqtt; print('MQTT: OK')" 2>/dev/null || log_error "Missing module: paho-mqtt"
python3 -c "import smbus; print('SMBus: OK')" 2>/dev/null || log_error "Missing module: smbus"
python3 -c "import PIL; print('PIL: OK')" 2>/dev/null || log_error "Missing module: PIL"

# Check I2C devices
log_info "Checking I2C devices..."
ls -la /dev/i2c-* 2>/dev/null || log_warning "No I2C devices found"

# Check Python file
log_info "Checking display_controller.py..."
if [ -f "/display_controller.py" ]; then
    log_info "File found: /display_controller.py"
    log_info "Size: $(stat -c%s /display_controller.py) bytes"
    # Python syntax test
    python3 -m py_compile /display_controller.py 2>/dev/null && log_info "Python syntax: OK" || log_error "Invalid Python syntax"
else
    log_error "/display_controller.py not found"
    exit 1
fi

log_info "Starting OLED controller..."
exec python3 -u /display_controller.py