# oled-i2c-screen-adapter/rootfs/etc/services.d/display_controller/run
#!/bin/sh
CONFIG_FILE="/data/options.json"

# Configuration MQTT
MQTT_BROKER=$(jq -r '.mqtt_broker' "$CONFIG_FILE")
MQTT_PORT=$(jq -r '.mqtt_port' "$CONFIG_FILE")
MQTT_USER=$(jq -r '.mqtt_user' "$CONFIG_FILE")
MQTT_PASSWORD=$(jq -r '.mqtt_password' "$CONFIG_FILE")

# Configuration I2C et écran
I2C_ADDRESS=$(jq -r '.i2c_address' "$CONFIG_FILE")
I2C_PORT=$(jq -r '.i2c_port' "$CONFIG_FILE")
DISPLAY_WIDTH=$(jq -r '.display_width' "$CONFIG_FILE")
DISPLAY_HEIGHT=$(jq -r '.display_height' "$CONFIG_FILE")
DISPLAY_TYPE=$(jq -r '.display_type' "$CONFIG_FILE")

# Configuration des paramètres d'affichage
REFRESH_INTERVAL=$(jq -r '.refresh_interval' "$CONFIG_FILE")
AUTO_BRIGHTNESS=$(jq -r '.auto_brightness' "$CONFIG_FILE")
DEFAULT_BRIGHTNESS=$(jq -r '.default_brightness' "$CONFIG_FILE")

# Configuration des topics MQTT
MQTT_TOPIC_TEXT=$(jq -r '.topics.text' "$CONFIG_FILE")
MQTT_TOPIC_COMMAND=$(jq -r '.topics.command' "$CONFIG_FILE")
MQTT_TOPIC_MODE=$(jq -r '.topics.mode' "$CONFIG_FILE")
MQTT_TOPIC_BRIGHTNESS=$(jq -r '.topics.brightness' "$CONFIG_FILE")
MQTT_TOPIC_REFRESH=$(jq -r '.topics.refresh' "$CONFIG_FILE")

# Configuration Tierce
QR_LINK=$(jq -r '.qr_link' "$CONFIG_FILE")

# Export des variables d'environnement
export MQTT_BROKER
export MQTT_PORT
export MQTT_USER
export MQTT_PASSWORD
export I2C_ADDRESS
export I2C_PORT
export DISPLAY_WIDTH
export DISPLAY_HEIGHT
export DISPLAY_TYPE
export REFRESH_INTERVAL
export AUTO_BRIGHTNESS
export DEFAULT_BRIGHTNESS
export MQTT_TOPIC_TEXT
export MQTT_TOPIC_COMMAND
export MQTT_TOPIC_MODE
export MQTT_TOPIC_BRIGHTNESS
export MQTT_TOPIC_REFRESH
export QR_LINK

echo "=== Configuration chargée ==="
echo "MQTT Broker: $MQTT_BROKER:$MQTT_PORT"
echo "MQTT User: $MQTT_USER"
echo "Configuration I2C: Port $I2C_PORT, Adresse $I2C_ADDRESS"
echo "Configuration écran: ${DISPLAY_WIDTH}x${DISPLAY_HEIGHT} ($DISPLAY_TYPE)"
echo "Refresh: ${REFRESH_INTERVAL}s, Brightness: $DEFAULT_BRIGHTNESS"
echo "Topics MQTT:"
echo "  - Text: $MQTT_TOPIC_TEXT"
echo "  - Command: $MQTT_TOPIC_COMMAND"
echo "  - Mode: $MQTT_TOPIC_MODE"
echo "  - Brightness: $MQTT_TOPIC_BRIGHTNESS"
echo "  - Refresh: $MQTT_TOPIC_REFRESH"
echo "=========================="

# Vérification des dépendances Python
echo "Vérification des dépendances Python..."
python3 -c "import os, sys; print(f'Python version: {sys.version}')"
python3 -c "import json; print('JSON: OK')" 2>/dev/null || echo "ERREUR: Module json manquant"
python3 -c "import paho.mqtt.client as mqtt; print('MQTT: OK')" 2>/dev/null || echo "ERREUR: Module paho-mqtt manquant"
python3 -c "import smbus; print('SMBus: OK')" 2>/dev/null || echo "ERREUR: Module smbus manquant"
python3 -c "import PIL; print('PIL: OK')" 2>/dev/null || echo "ERREUR: Module PIL manquant"

# Vérification des périphériques I2C
echo "Vérification des périphériques I2C..."
ls -la /dev/i2c-* 2>/dev/null || echo "ATTENTION: Périphériques I2C non trouvés"

# Vérification du fichier Python
echo "Vérification du fichier display_controller.py..."
if [ -f "/display_controller.py" ]; then
    echo "Fichier trouvé: /display_controller.py"
    echo "Taille: $(stat -c%s /display_controller.py) bytes"
    # Test de syntaxe Python
    python3 -m py_compile /display_controller.py 2>/dev/null && echo "Syntaxe Python: OK" || echo "ERREUR: Syntaxe Python invalide"
else
    echo "ERREUR: Fichier /display_controller.py non trouvé"
    exit 1
fi

echo "Lancement du controleur OLED..."
exec python3 /display_controller.py
